# gps-to-mqtt.service
[Unit]
Description=GPS to MQTT Bridge Service
Documentation=https://github.com/askrejans/gps-to-mqtt
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gps
Group=gps
WorkingDirectory=/
ExecStart=/usr/bin/gps-to-mqtt --config /etc/gps-to-mqtt/settings.toml

# Restart configuration
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

# Optional environment variable overrides (uncomment and adjust as needed):
# Environment=GPS_TO_MQTT_MQTT_HOST=localhost
# Environment=GPS_TO_MQTT_MQTT_PORT=1883
# Environment=GPS_TO_MQTT_LOG_LEVEL=info

# Resource limits
LimitNOFILE=65536
MemoryLimit=512M

# Environment
Environment="RUST_LOG=info"
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gps-to-mqtt

[Install]
WantedBy=multi-user.target
